# github-callback.service
[Unit]
Description=GitHub Callback Server für dennis AI Automation
Documentation=https://github.com/meinzeug/autodevai
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dennis
Group=dennis
WorkingDirectory=/home/dennis/autodevai
ExecStart=/usr/bin/node github-callback-server.js
ExecReload=/bin/kill -HUP $MAINPID

# Environment
Environment=NODE_ENV=production
Environment=PORT=19000
EnvironmentFile=-/home/dennis/autodevai/.env.dennis

# Process management
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/dennis/autodevai /tmp
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-callback-server

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target